# ansible/roles/vnc_server/templates/novnc@.service.j2
[Unit]
Description=noVNC Web VNC Client %i
After=network.target vnc@%i.service
Requires=vnc@%i.service

[Service]
Type=simple
User={{ ansible_user }}
# Use full path for websockify from system packages - no auth plugin needed
ExecStart=/usr/bin/websockify --web=/opt/novnc {{ 6080 + ('%i'|int) }} localhost:{{ 5900 + ('%i'|int) }}
# Add options for better performance and security
# --web-auth is available if you want to add HTTP basic auth later
Environment=PYTHONPATH=/opt/novnc
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target