# ansible/roles/vnc_server/templates/vnc.service.j2
[Unit]
Description=VNC Server for Remote Access
After=xvfb.service
Requires=xvfb.service

[Service]
User={{ ansible_user }}
Environment=DISPLAY=:1
ExecStart=/usr/local/bin/start-vnc-clipped.sh auto
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target