# ansible/roles/server/templates/kindle-automator.service.j2
[Unit]
Description=Kindle Automator Flask Server
After=network.target

[Service]
User=root
WorkingDirectory=/opt/kindle-automator
Environment="PATH=/opt/kindle-automator/venv/bin:/opt/android-sdk/platform-tools:/opt/android-sdk/emulator:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"
Environment="ANDROID_HOME=/opt/android-sdk"
Environment="ANDROID_SDK_ROOT=/opt/android-sdk"
Environment="FLASK_ENV=production"
Environment="ENVIRONMENT=prod"
Environment="PYTHONPATH=/opt/kindle-automator"
ExecStart=/opt/kindle-automator/venv/bin/python -m server.server
# Run graceful shutdown before stopping the service
ExecStop=/opt/kindle-automator/venv/bin/python /opt/kindle-automator/scripts/graceful_shutdown.py
# Add timeout settings for proper service termination
TimeoutStartSec=60s
TimeoutStopSec=30s
# Use SIGINT for graceful shutdown, then SIGTERM after timeout
KillMode=mixed
KillSignal=SIGINT
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 
